# extra/autostart/battery.desktop
[Desktop Entry]
Name=cbatticon
Type=Application
Exec=cbatticon