# extra/autostart/telegram.desktop
[Desktop Entry]
Name=Telgram
Type=Application
Exec=telegram-desktop -startintray
Icon=org.telegram.desktop